# scripts/kaivm-agent.service
[Unit]
Description=kaiVM Agent Web Service
After=network.target kaivm-capture.service
Wants=kaivm-capture.service

[Service]
Type=simple
User=kai
Group=kai
WorkingDirectory=/home/kai/kaivm
ExecStart=/home/kai/kaivm/.venv/bin/kaivm serve --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# Environment=GEMINI_API_KEY=... (set this in override or /etc/default/kaivm)

[Install]
WantedBy=multi-user.target